# etc/init.d/eudev-openrc
#!/sbin/openrc-run

description="Starts the eudev daemon to manage device nodes"

command="/usr/sbin/udevd"
supervisor=supervise-daemon

depend() {
    need dev procfs sysfs
    provide udev
}

start_pre() {
    if [ -w /sys/kernel/uevent_helper ]; then
        echo > /sys/kernel/uevent_helper
    fi

    udevadm info --cleanup-db >/dev/null 2>&1
}

start_post() {
    ebegin "Triggering udev events"
    udevadm trigger --type=subsystems --action=add
    udevadm trigger --type=devices --action=add
    eend $?

    ebegin "Waiting for udev to settle"
    if udevadm settle; then
        eend 0
    else
        eend 0 "udevadm settle timed out, continuing boot"
    fi
}

reload() {
    ebegin "Reloading udev rules"
    udevadm control --reload-rules
    eend $?
}
