# etc/init.d/eudev-rules-refresh
#!/sbin/openrc-run

description="Reload udev rules and retag ALL devices"

depend() {
    need udev localmount
    before greetd net
}

start() {
    ebegin "Reloading udev rules"
    if udevadm control --reload-rules; then
        eend 0
    else
        eend $? "Failed to reload udev rules"
        return 1
    fi

    ebegin "Retagging ALL THE BLOODY devices"
    if udevadm trigger --action=change; then
        eend 0
    else
        eend $? "Failed to retrigger udev for devices"
        return 1
    fi

    ebegin "Waiting for udev to settle"
    if udevadm settle --timeout=8; then
        eend 0
    else
        eend 0 "udevadm settle timed out, continuing"
    fi
}

stop() {
    return 0
}
